# backup.cron
# Run the backup script daily at 2 AM.
0 2 * * * /usr/local/bin/backup.sh >> /var/log/backup.log 2>&1
